# emacs-daemon.plist
<plist version="1.0">
 <dict>
   <key>Label</key>
   <string>emacs-daemon</string>
   <key>RunAtLoad</key>
   <true />
   <key>Program</key>
   <string>/Users/nicholas/.emacs-daemon.sh</string>
 </dict>
</plist>

       






